# animebot/animebot.service
[Unit]
Description=Anime Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/animebot
EnvironmentFile=/home/animebot/.env
ExecStart=/usr/bin/python3 /home/animebot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
